# openwrt/packages/luci-theme-thelinkbox/files/etc/uci-defaults/luci-theme-thelinkbox
#!/bin/sh
uci batch <<-EOF
	set luci.themes.thelinkbox=/luci-static/thelinkbox
	set luci.main.mediaurlbase=/luci-static/thelinkbox
        commit luci
EOF
	
